# bot.conf
[program:likestbot]
directory=/home/ubuntu/likestbot
command=python3 app.py

autostart=true
autorestart=true
enviroment=HOME="/home/ubuntu", USER="ubuntu"

stderr_logfile=/home/ubuntu/likestbot/logfile_err.log
stdout_logfile=/home/ubuntu/likestbot/logfile.log